# services/camera/camera-heartbeat.service
[Unit]
Description=Camera Heartbeat Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/camera_node
Environment=CN_CONFIG=/home/pi/camera_node/config.yaml
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/camera_node/venv/bin/python /home/pi/camera_node/heartbeat_client.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
