# services/camera/camera-heartbeat.service
[Unit]
Description=Camera Heartbeat Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/camera_node
# Use your actual config file (json in your case)
Environment=CN_CONFIG=/home/pi/camera_node/config.yaml
# Force unbuffered stdout/stderr so logs appear live in journal
Environment=PYTHONUNBUFFERED=1
# Use the venv interpreter so 'requests' etc. are available
ExecStart=/home/pi/camera_node/venv/bin/python /home/pi/camera_node/heartbeat_client.py
# Always restart so transient/bad-exit cases keep retrying
Restart=always
RestartSec=5
# Send all output to journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

